# supervisord-mailer.conf
[program:nullmailer]
command=/root/scripts/start-mailer.sh
numprocs=1
autostart=true
autorestart=true
startsecs=3